# config/systemd/system/marlin.timer
[Unit]
Description=Fetch and decode Marlin data from email every hour

[Timer]
OnCalendar=hourly
Unit=marlin-decode.service

[Install]