# config/systemd/user/decode-iridium-dump.path
[Unit]
Description=Dump Iridium data if any files are present in ~/Mail/attachments/
PartOf=decode-iridium.target

[Path]
DirectoryNotEmpty=%h/Mail/attachments/
Unit=decode-iridium-dump.service

[Install]
RequiredBy=decode-iridium.target
